# preview/supervisord.conf
# LeadCMS Live Preview Supervisor Configuration
# Manages multiple services in a single container:
# - nginx: Reverse proxy on port 80
# - livepreview: Development server on port 3000 (framework-agnostic)
# - sse-watcher: LeadCMS live content updates

[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
silent=false

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
redirect_stderr=false
priority=100

[program:livepreview]
# Framework-agnostic development server command
# Your project should implement: pnpm run livepreview
# Examples:
# - Next.js: "next dev"
# - Astro: "astro dev"
# - Gatsby: "gatsby develop"
# - Nuxt: "nuxt dev"
command=pnpm run livepreview
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
redirect_stderr=false
environment=NODE_ENV=development
priority=200

[program:sse-watcher]
# LeadCMS SSE watcher for live content updates
command=pnpm exec leadcms watch
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
redirect_stderr=false
environment=NODE_ENV=development
priority=300
# Only start if LeadCMS configuration is available
startsecs=10
